# scripts/cloudflared.service
[Unit]
Description=Cloudflare Tunnel
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/local/bin/cloudflared tunnel run career-advisor
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
